# deploy/cmms-updater.service
# CMMS host-side updater agent (oneshot).
#
# Performs the privileged half of the in-app "Update CMMS" button: when the
# webapp drops an update request into deploy-control/requests/, this service
# runs `git pull && docker compose up -d --build`.
#
# Install:
#   sudo cp deploy/cmms-updater.service /etc/systemd/system/
#   sudo cp deploy/cmms-updater.timer   /etc/systemd/system/
#   # Edit the paths/user below to match your host, then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now cmms-updater.timer
#
# See XCLOUD-DEPLOYMENT.md "One-click update button" for the full guide.

[Unit]
Description=CMMS host-side updater (git pull + docker compose rebuild)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot

# --- EDIT THESE FOR YOUR HOST ---
# The user must be able to run git pull and `docker compose` for this project.
User=deploy
# Absolute path to the cloned repo (the docker compose project root).
WorkingDirectory=/opt/cmms
Environment=REPO_DIR=/opt/cmms
Environment=DEPLOY_CONTROL_DIR=/opt/cmms/deploy-control
Environment=COMPOSE_CMD=docker compose
# Optionally pin the branch (defaults to the repo's current branch):
# Environment=GIT_BRANCH=main
# --------------------------------

ExecStart=/opt/cmms/scripts/host-updater.sh once
